# com.buildone.sync-qbo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.buildone.sync-qbo</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/build.one/.venv/bin/python3</string>
        <string>-m</string>
        <string>scripts.sync_qbo_company_info</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/build.one</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/path/to/build.one/logs/sync-qbo.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/build.one/logs/sync-qbo-error.log</string>
</dict>
</plist>
